# systemd/.config/systemd/user/laptop-power-apply.service
[Unit]
Description=Apply laptop AC/battery power policy
StartLimitIntervalSec=0

[Service]
Type=oneshot
KillMode=process
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=%h/.local/bin/laptop-power-apply
