# systemd/.config/systemd/user/laptop-power-apply.timer
[Unit]
Description=Poll laptop AC/battery state

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=laptop-power-apply.service

[Install]
WantedBy=timers.target
